# vnc.service
[Unit]
Description=OneLab Server

[Service]
Path=your_path
ExecStart=your_bin
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
