# vnc.service
[Unit]
Description=OneLab VNC

[Service]
ExecStart=/usr/bin/x11vnc --loop
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
